# projects/inner-light/src/scripts/innerlightdrive.crontab
# min hour day-month month day-week cmd
* * * * * $HOME/bin/innerlightdrive-monitor >> $HOME/data/innerlightdrive-monitor.log 2>> $HOME/data/innerlightdrive-monitor.err
